# recipes/middleware/kibana/templates/kibana.conf
start on (local-filesystems and net-device-up IFACE!=lo)
stop on shutdown

respawn
setuid root
cd /var/www/kibana

exec ./bin/kibana >> /var/log/kibana 2>&1
